# systemd/codex-langfuse-watch.service
[Unit]
Description=Export completed Codex turns to Langfuse
After=default.target

[Service]
Type=simple
ExecStart=%h/.codex/bin/codex-langfuse-exporter --watch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
